# etc/cron.hourly/backup
#!/bin/sh

/root/backup/export.sh

test -x /usr/sbin/logrotate || exit 0
/usr/sbin/logrotate /root/backup/rotate.conf
