# etc/cron.hourly/backup
#!/bin/sh

if test -x /usr/sbin/logrotate
then
  /usr/sbin/logrotate /root/backup/rotate.conf
fi

/root/backup/export.sh
